# systemd/tradingapp-market-hours.service
[Unit]
Description=Trading App Market Data Services
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/tradingapp
# Force restart services to pick up fresh tokens
ExecStartPre=/usr/bin/docker-compose stop ingestion depth-collector signal-generator worker-1 worker-2 worker-3
ExecStartPre=/usr/bin/docker-compose rm -f ingestion depth-collector signal-generator worker-1 worker-2 worker-3
ExecStartPre=/usr/bin/docker-compose build ingestion depth-collector signal-generator worker-1 worker-2 worker-3
ExecStart=/usr/bin/docker-compose up -d ingestion depth-collector signal-generator worker-1 worker-2 worker-3
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
