# systemd/tradingapp-start.timer
[Unit]
Description=Start Trading App at 9:00 AM IST

[Timer]
# TEMPORARY TEST: Run at 19:10 CET (18:10 UTC) for testing
OnCalendar=*-*-* 18:10:00
Persistent=true
Unit=tradingapp-market-hours.service

[Install]
WantedBy=timers.target
