# systemd/tradingapp-start.timer
[Unit]
Description=Start Trading App at 9:00 AM IST

[Timer]
# TEMPORARY TEST: Run at 11:25 PM IST (17:55 UTC) for testing
OnCalendar=*-*-* 17:55:00
Persistent=true
Unit=tradingapp-market-hours.service

[Install]
WantedBy=timers.target
