# systemd/tradingapp-start.timer
[Unit]
Description=Start Trading App at 9:00 AM IST

[Timer]
# TEMPORARY TEST: Run at 11:30 PM IST (18:00 UTC) for testing
OnCalendar=*-*-* 18:00:00
Persistent=true
Unit=tradingapp-market-hours.service

[Install]
WantedBy=timers.target
