# systemd/tradingapp-stop.timer
[Unit]
Description=Stop Trading App at 3:40 PM IST

[Timer]
# 3:40 PM IST = 10:10 UTC (Monday-Friday)
OnCalendar=Mon-Fri *-*-* 10:10:00 UTC
Persistent=true
Unit=tradingapp-market-hours-stop-pm2.service

[Install]
WantedBy=timers.target
