# systemd/tradingapp-stop.timer
[Unit]
Description=Stop Trading App at 3:40 PM IST

[Timer]
# TEMPORARY TEST: Run at 18:20 UTC for testing
OnCalendar=*-*-* 18:20:00 UTC
Persistent=true
Unit=tradingapp-market-hours-stop.service

[Install]
WantedBy=timers.target
