# systemd/tradingapp-market-hours-stop.service
[Unit]
Description=Stop Trading App Market Data Services
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/tradingapp
ExecStart=/usr/bin/docker-compose stop ingestion depth-collector signal-generator worker-1 worker-2 worker-3
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
