# systemd/tradingapp-start.timer
[Unit]
Description=Start Trading App at 9:00 AM IST

[Timer]
# 9:00 AM IST = 03:30 UTC (Monday-Friday)
OnCalendar=Mon-Fri *-*-* 03:30:00 UTC
Persistent=true
Unit=tradingapp-market-hours-pm2.service

[Install]
WantedBy=timers.target
